# files/profile.d/history.sh
#!/bin/bash
# /etc/profile.d/history.sh
# hitory に時刻情報を追記(全ユーザ)

HISTTIMEFORMAT="%F %T "
